# leaf-by-leaf.service
[Unit]
Description=Leaf by Leaf — personal voice journaling app
After=network.target

[Service]
User=will
WorkingDirectory=/home/will/apps/leaf-by-leaf
ExecStart=/home/will/apps/leaf-by-leaf/venv/bin/gunicorn -w 1 -b 0.0.0.0:8005 --timeout 120 app:app
Restart=on-failure
RestartSec=5

EnvironmentFile=/home/will/.env/leaf-by-leaf

[Install]
WantedBy=multi-user.target
